# src/conf/systemd/systemd-firstboot.service.d/prompt.conf
[Service]
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/console
