# src/conf/systemd/pacman-auto-upgrade.service
[Unit]
Description=Daily pacman -Syu auto-upgrade (snapper snapshots via snap-pac)

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Syu --noconfirm
